# webui/launchd/com.example.ledger-webui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS 자동 시작(launchd) 템플릿.
  사용법:
    1) 아래 __LEDGER_DIR__ 을 본인의 설치 경로(예: /Users/<사용자명>/Documents/Ledger)로 모두 바꾼다.
    2) 파일명을 com.<본인id>.ledger-webui.plist 로 바꿔 ~/Library/LaunchAgents/ 에 복사한다.
    3) launchctl load ~/Library/LaunchAgents/com.<본인id>.ledger-webui.plist
  자동 시작이 필요 없으면 이 파일은 무시하고 webui/run.command 더블클릭으로 실행해도 된다.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.ledger-webui</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>PORT="${PORT:-5001}"; if lsof -i :"$PORT" -sTCP:LISTEN &gt;/dev/null 2&gt;&amp;1; then echo "Port $PORT in use; skip."; exit 0; fi; exec __LEDGER_DIR__/webui/run.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>5001</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/hledger-web.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hledger-web.launchd.log</string>
</dict>
</plist>
